# config/supervisormazey/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/web/supervisor_nginx.log
stderr_logfile=/var/log/web/supervisor_nginx_err.log
priority=2
